# auth-server/systemd/update-player-elo-daily.service
[Unit]
Description=Update Carcassonne GG player Elo for all players and refresh planned ratings
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/carcassonne-gg/auth-server
LogsDirectory=carcassonne
StandardOutput=append:/var/log/carcassonne/update-player-elo-daily.log
StandardError=append:/var/log/carcassonne/update-player-elo-daily.log
ExecStart=/home/carcassonne-gg/auth-server/scripts/run_update_player_elo_daily.sh --db-path /home/carcassonne-gg/auth-server/data/auth.sqlite
TimeoutStartSec=30min

[Install]
WantedBy=multi-user.target
